# Bitnami/AcctAggregation/crontab
* * * * *   /usr/bin/execute /aggregate.sh
