# Bitnami/AcctAggregation/crontab
@reboot /config.slaves.sh
*/5 * * * *   /aggregate.sh
