# Bitnami/AcctAggregation/crontab
@reboot /config.slaves.sh
*/2 * * * *   /aggregate.sh
